# dockervault_ghostunnel_server_1_6_1_go_1_19_x64/files/ghostunnel_server.conf
[program:ghostunnel_server]
command=/var/opt/ghostunnel/ghostunnel server --shutdown-timeout=2m --connect-timeout=10s --unsafe-target --listen 0.0.0.0:%(ENV_LISTEN_PORT)s --target %(ENV_TARGET_HOST_AND_PORT)s --keystore /etc/ghostunnel_server/ghostunnel_server.pfx --cacert /etc/ghostunnel_server/ghostunnel_ca.cer --allow-cn ghostunnel_Client
user=root
stopsignal=KILL
autostart=true
autorestart=true
startsecs=30
exitcodes=0
stopwaitsec=0
stopasgroup=true
startretries=2147483647
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_syslog=true
stderr_syslog=true

